# templates/profile.d/000.homebrew.sh
if [[ "$OSTYPE" == "linux-gnu"* ]]; then
	eval $(/home/linuxbrew/.linuxbrew/bin/brew shellenv)
fi
