# templates/profile.d/002.docker.sh
#/bin/bash

set -e

alias dc=docker-compose
