# templates/profile.d/002.autojump.sh
function j() {
  . /usr/share/autojump/autojump.sh
  j "$@"
}
