# templates/profile.d/123.klog.sh
TRACK_DIR=$HOME/Documents/time-tracking
alias kreport='klog report $TRACK_DIR/$(ls $TRACK_DIR | tail -n1)'
